# templates/kojid.ulimits.service.j2
[Service]
LimitNOFILE={{ kojid_ulimits_nofiles }}
LimitNOFILESoft={{ kojid_ulimits_nofiles }}
